# deploy/almalinux/api.service
[Unit]
Description=Sis Teknik API (Rust Axum)
After=network.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/opt/sis-teknik
ExecStart=/opt/sis-teknik/backend/target/release/api
Restart=always
RestartSec=3
Environment=API_PORT=3000
Environment=DATABASE_URL=postgresql://sis_teknik_app:SisTeknikApp_2026!Pg@127.0.0.1:5432/sis_teknik
Environment=FRONTEND_URL=https://tamir.sis-teknik.com.tr
Environment=TURNSTILE_SECRET_KEY=0x4AAAAAAChToBHieN0cwGcSSMhurhQHm64

[Install]
WantedBy=multi-user.target
